# docker/bot/cron.d/python-job
0 6 * * * /bin/bash /app/posttoday.sh >> /var/log/cron.log 2>&1

*/15 * * * * /bin/bash /app/postnew.sh >> /var/log/cron.log 2>&1
